# systemd/ai-energy-update-tariffs.service
[Service]
User=saltspork
Type=oneshot
Nice=19
WorkingDirectory=/home/saltspork/src/ai-energy-forecast-slop
ExecStart=/bin/bash -c 'source .venv/bin/activate && ./forecast.py update-tariffs && ./smooth_tariffs.py && ./forecast.py backfill-actuals && ./forecast.py update-adjusters'
