# systemd/ai-energy-predict.service
[Unit]
Description=Run AI energy predictions and publish to Home Assistant

[Service]
Type=oneshot
WorkingDirectory=/home/saltspork/src/ai-energy-forecast-slop
EnvironmentFile=/home/saltspork/src/ai-energy-forecast-slop/.env
ExecStart=/bin/bash -c 'source .venv/bin/activate && ./forecast.py predict-all --dynamic-handoff --publish-hass --publish-covariates && curl -fsS --retry 3 ${HC_PREDICT_URL}'
