# systemd/ai-energy-pd7day.service
[Unit]
Description=Ingest latest AEMO PD7Day price forecast into InfluxDB

[Service]
Type=oneshot
WorkingDirectory=/home/saltspork/src/ai-energy-forecast-slop
ExecStart=/bin/bash -c 'source .venv/bin/activate && python ingest/ingest-pd7day.py --fetch'
