# systemd/ai-energy-sevendayoutlook.service
[Unit]
Description=Ingest latest AEMO SEVENDAYOUTLOOK_FULL demand forecast into InfluxDB

[Service]
Type=oneshot
WorkingDirectory=/home/saltspork/src/ai-energy-forecast-slop
ExecStart=/bin/bash -c 'source .venv/bin/activate && python ingest/ingest-sevendayoutlook.py --fetch'
